# mkosi.profiles/runtime/mkosi.extra/usr/lib/systemd/system/sshd-keygen@.service.d/10-katl-persistent-host-key.conf
[Unit]
# Katl owns the persistent host key instead of writing host keys to immutable /etc.
ConditionPathExists=/var/lib/katl/ssh/enable-distribution-host-key-generator
