# mkosi.profiles/installer-image/mkosi.extra/usr/lib/systemd/system/getty@tty1.service.d/10-katl-console.conf
[Unit]
ConditionPathExists=!/usr/bin/katl-console
