# mkosi.profiles/runtime/mkosi.extra/usr/lib/systemd/system/sshd.service.d/10-katl-host-keys.conf
[Unit]
Wants=katl-ssh-host-keys.service
After=katl-ssh-host-keys.service
